# docker-compose/web/app.conf
[supervisord]
nodaemon=true

[program:spark-master]
command=node /app/index.js
