# cluster/spark/master-web.conf
[supervisord]
nodaemon=true

[program:web]
command=node /app/index.js
